# copy/var/tmp/dcm4chee-2.17.3-mysql/bin/dcm4chee-manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='network/dcm4chee' type='service' version='0'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='net' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/loopback'/>
    </dependency>
    <exec_method name='start' type='method' exec='/export/dcm4chee/bin/svc-dcm4chee start' timeout_seconds='18446744073709551615'/>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='18446744073709551615'/>
    <exec_method name='restart' type='method' exec='/export/dcm4chee/bin/svc-dcm4chee restart' timeout_seconds='18446744073709551615'/>
    <stability value='Unstable'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>dcm4chee</loctext>
      </common_name>
      <documentation>
        <manpage title='dcm4chee' section='1M' manpath='/usr/share/man'/>
      </documentation>
    </template>
  </service>
</service_bundle>
